# proutgpt.service
[Unit]
Description=ProutGPT Backend (OpenRouter Proxy)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/proutgpt-backend
# Load the API key from /etc/proutgpt.env (KEY=value format, no 'export')
# Create it with: echo "OPENROUTER_API_KEY=sk-or-v1-..." | sudo tee /etc/proutgpt.env && sudo chmod 600 /etc/proutgpt.env
EnvironmentFile=/etc/proutgpt.env
ExecStart=/home/ubuntu/openrouter-env/bin/python /home/ubuntu/proutgpt-backend/openrouter_proxy.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/ubuntu/openrouter.log
StandardError=append:/home/ubuntu/openrouter.log

[Install]
WantedBy=multi-user.target
